# systemd/vds-watchdog.service
[Unit]
Description=VDS ZeroTier Watchdog
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/my-vds
Environment=PYTHONPATH=/opt/my-vds
ExecStart=/opt/my-vds/venv/bin/python -m src.main zerotier watchdog
StandardOutput=journal
StandardError=journal
IOSchedulingClass=idle
CPUWeight=1
