# systemd/vds-watchdog.timer
[Unit]
Description=VDS ZeroTier Watchdog Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=30s

[Install]
WantedBy=timers.target
